# _config/cinemanager_django.socket
[Unit]
Description=gunicorn socket for cinemanager app

[Socket]
ListenStream=/run/cinemanager_django.sock

[Install]
WantedBy=sockets.target
