# _config/sites_django.socket
[Unit]
Description=gunicorn socket for sites django app

[Socket]
ListenStream=/run/sites_django.sock

[Install]
WantedBy=sockets.target
