# setup-conf/systemd/bounces.service
[Unit]
Description=Checks for bounced newsletter users

[Service]
Type=oneshot
ExecStart=/bin/bash /var/www/Cinemanager/_cron/bounces.sh
