# setup-conf/systemd/newsletter.service
[Unit]
Description=Sends scheduled newsletters

[Service]
Type=oneshot
ExecStart=/bin/bash /var/www/Cinemanager/_cron/newsletter.sh

